# launchd/com.denderson.ademir.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>ADEMIR_ACTIVE_END</key>
		<string>18</string>
		<key>ADEMIR_ACTIVE_START</key>
		<string>9</string>
		<key>ADEMIR_DAILY_DM_LIMIT</key>
		<string>30</string>
		<key>ADEMIR_DELAY_BETWEEN_LEADS_SEC</key>
		<string>0</string>
		<key>ADEMIR_DELAY_JITTER_SEC</key>
		<string>0</string>
		<key>ADEMIR_DISABLE_CRON</key>
		<string>true</string>
		<key>ADEMIR_DRY_RUN</key>
		<string>true</string>
		<key>ADEMIR_HIKER_FOLLOWERS_PAGES</key>
		<string>1000</string>
		<key>ADEMIR_PER_RUN_CAP</key>
		<string>500</string>
		<key>ADEMIR_PORT</key>
		<string>9100</string>
		<key>ADEMIR_SKIP_DM_ON_DISCOVERY</key>
		<string>true</string>
		<key>ADEMIR_TOKEN</key>
		<string>{{ADEMIR_TOKEN}}</string>
		<key>CLAUDE_CODE_OAUTH_TOKEN</key>
		<string>{{CLAUDE_CODE_OAUTH_TOKEN}}</string>
		<key>HIKERAPI_KEY</key>
		<string>{{HIKERAPI_KEY}}</string>
		<key>OPENAI_API_KEY</key>
		<string>{{OPENAI_API_KEY}}</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>PG_DB</key>
		<string>naia_memory</string>
		<key>PG_HOST</key>
		<string>127.0.0.1</string>
		<key>PG_PASS</key>
		<string>{{POSTGRES_PASSWORD}}</string>
		<key>PG_PORT</key>
		<string>15432</string>
		<key>PG_USER</key>
		<string>n8n</string>
		<key>TANDEM_TOKEN</key>
		<string>{{TANDEM_TOKEN}}</string>
		<key>TANDEM_URL</key>
		<string>http://localhost:8765</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.{{DONO_SLUG}}.ademir</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/python3</string>
		<string>/Users/naiarodrigues/naia-agent/scripts/ademir.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/naiarodrigues/naia-agent/logs/ademir.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/naiarodrigues/naia-agent/logs/ademir.out.log</string>
</dict>
</plist>
